# .config/systemd/user/vault-backup.service
[Unit]
Description=Backup keepassxc vault

[Service]
Type=oneshot
ExecStart=%h/.scripts/backup-vault.sh

[Install]
WantedBy=default.target
