# .config/systemd/user/vault-backup.timer
[Unit]
Description=Backup keepassxc vault

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
